# templates/default/kibana.service.erb
[Unit]
Description= Flashcards kibana instance

[Service]
User=<%= node["kibana"]["user"] %>
Group=<%= node["kibana"]["group"] %>
WorkingDirectory=/opt/kibana
PIDFile=/var/run/kibana.pid
SyslogIdentifier=flashcards-kibana
Restart=always
ExecStart=/usr/local/bin/kibana

[Install]
WantedBy=multi-user.target