# templates/default/kibana.service.erb
[Unit]
Description= Flashcards Kibana 4 instance

[Service]
Type=simple
User=<%= node["kibana"]["user"] %>
Group=<%= node["kibana"]["group"] %>
WorkingDirectory=/opt/kibana-<%= node["kibana"]["version"] %>
SyslogIdentifier=flashcards-kibana
Restart=always
ExecStart=/opt/kibana-<%= node["kibana"]["version"] %>/bin/kibana

[Install]
WantedBy=multi-user.target
