# templates/default/kibana.service.erb
[Unit]
Description= Flashcards Kibana 4 instance

[Service]
User=<%= node["kibana"]["user"] %>
Group=<%= node["kibana"]["group"] %>
Environment=NODE_ENV=<%= node.environment %>
WorkingDirectory=/opt/kibana-<%= node["kibana"]["version"] %>
PIDFile=/var/run/flashcards-kibana.pid
SyslogIdentifier=flashcards-kibana
Restart=always
ExecStart=/opt/kibana-<%= node["kibana"]["version"] %>/bin/kibana

[Install]
WantedBy=multi-user.target
